# systemd/burrow.service
[Unit]
Description=Burrow
After=burrow.socket

[Service]
ExecStart=/usr/bin/burrow daemon

[Install]
WantedBy=multi-user.target
