# systemd/burrow.socket
[Unit]
Description=Burrow Socket

[Socket]
ListenStream=/run/burrow.sock

[Install]
WantedBy=sockets.target
